# services/midi-connector.service
[Unit]
Description=Midi Connector Web UI
After=network.target sound.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/apps/Midi-Connector
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
